# init/gstd.service.in
[Unit]
Description=GstD systemd service to generate folders for run-time variable data and log files.

[Service]
Type=forking
ExecStart=@prefix@/libexec/gstd/gstd start
ExecStop=@prefix@/libexec/gstd/gstd stop

[Install]
WantedBy=multi-user.target