# init/gstd.service.in
[Unit]
Description=GStreamer Daemon
PartOf=graphical-session.target

[Service]
Restart=on-failure
RestartSec=5s
ExecStart=@prefix@/bin/gstd

[Install]
WantedBy=graphical-session.target
