# init/gstd.service.in
[Unit]
Description=GstD systemd service to generate folders for run-time variable data and log files.
PartOf=graphical-session.target

[Service]
Type=forking
ExecStart=@PREFIX@/bin/gstd -e
ExecStop=@PREFIX@/bin/gstd -e -k

[Install]
WantedBy=graphical-session.target
