# init/gstd.service.in
[Unit]
Description=Systemd service to launch GstD
PartOf=graphical-session.target

[Service]
Type=forking
ExecStart=@prefix@/bin/gstd -e
ExecStop=@prefix@/bin/gstd -k

[Install]
WantedBy=graphical-session.target
