# init/gstd.in
#!/bin/sh

### BEGIN INIT INFO
# Provides: GstD
# Short-Description: start and stop GstD
# Description: Gstreamer Daemon is a framework for
#              controlling audio and video streaming
### END INIT INFO

manage_dir()
{
    dir="$1";
    mode="$2";
    user="$3";

    #Create directory if not created yet
    mkdir -m "$mode" -p "$dir";
    if [ $? -ne 0 ]; then
        echo "GstD init: Error, cannot create directory $dir" 1>&2;
        exit 1
    fi

    #Set the owner
    chown "$user" "$dir";
    if [ $? -ne 0 ]; then
        echo "GstD init: Error, cannot change $dir owner to $user:$group" 1>&2;
        exit 1
    fi
}

prefix=@prefix@
gstdrunstatedir=@GSTD_RUN_STATE_DIR@
gstdlogstatedir=@GSTD_LOG_STATE_DIR@

case "$1" in
    start)
        directories=("$gstdrunstatedir" "$gstdlogstatedir")

        #Manage pid and log directories
        for d in ${directories[@]}; do
            manage_dir "$d" 777 @GSTD_STATE_DIRS_OWNER@
        done

        #Start gstd
        @prefix@/bin/gstd -e;
        ;;
    stop)
        @prefix@/bin/gstd -k -e;
        ;;
    *)
        echo "$0 {start|stop}"
esac
